# software/autostart.service
[Unit]
Description=ActionPi

[Service]
ExecStart=/bin/bash -c "cd /home/action && python action.py"
Type=simple

[Install]
WantedBy=default.target
